# evilwoods-server/ansible/roles/searxng/templates/searxng-valkey.container.jinja
[Unit]
Description=Valkey Container Service

[Container]
AutoUpdate=registry
Image=docker.io/valkey/valkey:8.1-alpine
Volume=valkey_data.volume:/data
Pod=searxng.pod

[Service]
Restart=on-failure
