# bootc/evilwoods-server/context/config/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/50-trigger-time.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=
OnCalendar=Sat, 07:00
