# evilwoods-server/ansible/roles/fusion/templates/fusion.container.j2
[Unit]
Description={{ fusion_app_name }} Container Service

[Container]
AutoUpdate=registr
Image={{ fusion_image }}
Volume={{ fusion_data_volume }}:/data:Z,U
Network={{ evilwoods_container_network }}.network

[Service]
Restart=always

[Install]
WantedBy=default.target
