# bootc/evilwoods-server/config/usr/etc/containers/systemd/users/container-runner/caddy.container
[Unit]
Description=Caddy container sevice

[Container]
Image=docker.io/library/caddy:2.9

[Install]
WantedBy=default.target
