# bootc/evilwoods-server/context/config/usr/etc/profile.d/evilwoods.sh
PATH=/usr/local/share/xonsh-env/xbin:$PATH
export PATH
