# ansible/roles/common/templates/common-restic-backup.timer.j2
[Unit]
Description=Restic backup timer

[Timer]
OnCalendar=*-*-* 04:00:00

[Install]
WantedBy={{ evilwoods_service_target }}
