# ansible/roles/common/templates/common-restic-backup.service.j2
[Unit]
Description=Restic backup

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '{{ common_restic_backup_script }}'
