# ansible/roles/common/templates/evilcloud.target.j2
[Unit]
Description=Evilwoods Cloud Target

[Install]
WantedBy=default.target
