# evilwoods-server/ansible/roles/nextcloud/templates/nextcloud-valkey.container.j2
[Unit]
Description={{ nextcloud_app_name }} valkey Container/Service

[Container]
Image=docker.io/valkey/valkey:8-alpine
Volume=nextcloud-valkey-data.volume:/data
Pod=nextcloud.pod
